# launchd/com.cascade.lobsterhook.dispatcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.cascade.lobsterhook.dispatcher</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/bash</string>
      <string>__REPO_ROOT__/scripts/run-dispatcher.sh</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
      <key>LOBSTERHOOK_CONFIG</key>
      <string>__REPO_ROOT__/lobsterhook.toml</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/data/logs/dispatcher.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/data/logs/dispatcher.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>
  </dict>
</plist>
